# edge/bootstrap/ngrokssh.service
[Unit]
Description=Example systemd service.

[Service]
Type=simple
ExecStart=/opt/ngrok tcp 22

[Install]
WantedBy=multi-user.target